# configs/dotfile/autostart/IBus.desktop
[Desktop Entry]
Name=ibus
Icon=ibus
Exec=ibus start
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
